# guest/incus-gh-runner-guest.path
[Unit]
Description=Wait for the Incus GitHub runner runtime payload

[Path]
PathExists=/run/incus-gh-runner/payload.ready
Unit=incus-gh-runner-guest.service

[Install]
WantedBy=multi-user.target
